# hass/python-wyoming-faster-whisper/wyoming-faster-whisper.service
[Unit]
Description=Wyoming Server for Faster Whisper
Wants=network.target

[Service]
Restart=on-failure
User=wyoming-faster-whisper
Group=wyoming-faster-whisper
DynamicUser=true
CacheDirectory=wyoming-faster-whisper
Environment="HF_HOME=/var/cache/wyoming-faster-whisper/huggingface"
StateDirectory=wyoming-faster-whisper
WorkingDirectory=/var/lib/wyoming-faster-whisper

EnvironmentFile=/etc/conf.d/wyoming-faster-whisper
ExecStart=wyoming-faster-whisper \
    --device $DEVICE \
    --model $MODEL \
    --uri $URI \
    --beam-size $BEAM_SIZE \
    --language $LANGUAGE \
    --data-dir /var/lib/wyoming-faster-whisper/data \
    --download-dir /var/lib/wyoming-faster-whisper/data

[Install]
WantedBy=multi-user.target