# hass/home-assistant-core/home-assistant.service
[Unit]
Description=Home assistant
After=network.target

[Service]
Restart=on-failure
User=hass
Group=hass
LogsDirectory=hass
StateDirectory=hass
WorkingDirectory=/var/lib/hass
Environment=HOME=/var/lib/hass
Environment=PATH=/usr/bin:/var/lib/hass/.venv/bin

ExecStartPre=/usr/bin/python -m venv --upgrade /var/lib/hass/.venv
ExecStartPre=/var/lib/hass/.venv/bin/python -m pip install wheel
ExecStartPre=/var/lib/hass/.venv/bin/python -m pip install @WHEEL@
ExecStart=/var/lib/hass/.venv/bin/hass \
  --config /var/lib/hass/ \
  --log-file /var/log/hass/home-assistant.log \
  --log-rotate-days 1
RestartForceExitStatus=100
TimeoutSec=15min

AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
LockPersonality=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true

[Install]
WantedBy=multi-user.target
